# hertz_feed.timer
[Unit]
Description=30 second timer for the hertz_feed service

[Timer]
OnUnitActiveSec=30s
OnBootSec=30s

[Install]
WantedBy=timers.target